# deploy/my_information.conf
[program:myinformation]
command = /webapps/My_information/.env/bin/gunicorn_start               ; Command to start app
user = webapp                      	                                    ; User to run as
stdout_logfile = /webapps/My_information/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                  ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                         ;
